# systemd/user/tmux.service
[Unit]
Description=tmux session for user

[Service]
Type=forking
ExecStart=/usr/bin/tmux new-session -s default -d
ExecStop=/usr/bin/tmux kill-session -t default
EnvironmentFile=/home/akisarou/dotfiles/zsh/exports.sh
[Install]
WantedBy=default.target
; After=graphical-session.target
